# supervisord.conf
[supervisord]
nodaemon=true

[program:crond]
command=/usr/local/bin/update_status.sh

[program:apache2]
command=/usr/local/bin/apache2-foreground
