# supervisord.conf
[supervisord]
nodaemon=true

[program:crond]
environment=UPDATE_INTERVAL=%(ENV_UPDATE_INTERVAL)s
command=/usr/local/bin/update_status.sh

[program:apache2]
command=/usr/local/bin/apache2-foreground
